# sysdrv/tools/board/debian/overlay/etc/update-motd.d/11-help-text
#!/bin/sh

printf "\n"
printf " * Documentation:  https://wiki.luckfox.com\n"
printf " * Management:     https://www.luckfox.com\n"
